# conf/supervisor/conf.d/django_uwsgi.ini
[program:lunwen_api]
command=/opt/lunwen/venv_py3/bin/uwsgi --ini /opt/lunwen/conf/django/uwsgi.ini
autostart=true
autorestart=true
user=root
startsecs=10
redirect_stderr=true
stdout_logfile=/var/log/lunwen/api.log
stopasgroup=true
stopsignal=QUIT


[program:nginx]
command=/usr/local/nginx/sbin/nginx -c /usr/local/nginx/conf/nginx.conf -g 'daemon off;'
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/lunwen/api_nginx.log
stopasgroup=true
stopsignal=QUIT
priority=1
